# etc/cron.weekly/letsencrypt-renew
#!/bin/bash
/usr/local/sbin/certbot-auto renew >> /var/log/le-renew.log
